# deploy/tfc-app.service
# systemd unit for the TFC app.
# Install to /etc/systemd/system/tfc-app.service on the EC2 host, then:
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now tfc-app
#
# Runs on 127.0.0.1:5001 so it sits beside the existing Schnurr app
# (which uses :5000). nginx reverse-proxies tfc.schnurr-app.com to this port.

[Unit]
Description=TFC App (Flask/gunicorn)
After=network.target

[Service]
Type=simple
# Match this to the user that owns /opt/tfc-app (often "ubuntu" on Ubuntu AMIs).
User=ubuntu
WorkingDirectory=/opt/tfc-app
Environment=PORT=5001
# gunicorn lives in the project venv created during deploy (see DEPLOY.md).
ExecStart=/opt/tfc-app/.venv/bin/gunicorn --workers 2 --bind 127.0.0.1:5001 app:app
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
